# teoria/03-container/02-runtime/diego-ceja/nginx.container
[Container]
Image=docker.io/library/nginx:1.28
PublishPort=33412:80
Volume=ngx-html:/usr/share/nginx/html
