# teoria/03-container/02-runtime/luis-salazar/quadlet/nginx.container

[Container]
Image=nginx
PublishPort=8080:80
Volume=%h/html:/usr/share/nginx/html:ro
