# teoria/02-iac/02-network-debugging/luis-salazar/recon.service
[Unit]
Description=Monitor y recuperación de interfaz de red

[Service]
Type=simple
ExecStart=/home/debian/renet.sh
Restart=Always
[Install]
WantedBy=multi-user.target