# teoria/02-iac/02-network-debugging/javier-arias/roles/watchdog/files/watchdog.timer
[Unit]
Description=Run network watchdog every 5 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=5s
Unit=watchdog.service

[Install]
WantedBy=timers.target
