# teoria/02-iac/02-network-debugging/diego-ceja/monitor.service
# /etc/systemd/system/monitor.service
# Systemd service file that executes the monitoring script

[Unit]
Description=Connectivity Monitor Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c '/root/monitor.sh'
User=root
Group=root

[Install]
WantedBy=multi-user.target

