# teoria/02-iac/02-network-debugging/mauricio-morales/files/net-watchdog.timer
[Unit]
Description=Network Interface Restart Timer

[Timer]
OnCalendar=*-*-* *:*:00/30
AccuracySec=1s

[Install]
WantedBy=timers.target
