# teoria/02-iac/02-network-debugging/bug.service
# /etc/systemd/system/bug.service
# Systemd service file that randomly cuts the internet every (10-15) minutes

[Unit]
Description=Random Internet Disconnection Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'while true; do sleep $((RANDOM % 300 + 600)); IFACE=$(ip route | awk "/default/ {print \\$5; exit}"); ip link set "$IFACE" down; sleep $((RANDOM % 60 + 30)); ip link set "$IFACE" up; done'
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target

